# opendeltaiso/airootfs/etc/profile.d/setup.sh
sudo sh /etc/cinnamon-setup.sh
sudo sh /etc/flatpak-setup.sh